# crontab
# fire incremental updates (3am EST)
24 18 * * 1-5 bash -l -c 'cd /app && python run.py -p'
