# crontab
# fire incremental updates (3am EST)
3,13,23,33,43,53 * * * 0-6 bash -l -c 'cd /app && python run.py'
